# iaas/provisioning/common/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
#Environment="HTTP_PROXY=http://your-proxy:3128"
#Environment="HTTPS_PROXY=http://your-proxy:3128"
#Environment="NO_PROXY=localhost,127.0.0.1,.sock"
