# iaas/config/etc/profile.d/custom-env-vars

# =====================================================
# You can customize your env vars on the iaas here
# =====================================================

# Custom env variables
export COMPOSE_PROJECT_NAME="app"
