# iaas/config/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=PUT_YOUR_PROXY_HERE"
Environment="HTTPS_PROXY=PUT_YOUR_PROXY_HERE"
Environment="NO_PROXY=localhost,127.0.0.1,.sock"
