# iaas/provisioning/common/etc/cron.daily/docker-prune
#!/bin/bash

#
# Big cleaning to be sure that the system is clean at any time
#
docker system prune -f
